# deployments/launchers/base/src/main/bin/init/aprox.service
# 1. Adjust /opt/aprox/bin to match the deployment location
# 2. Put in /etc/systemd/system/
# 3. Execute `systemctl enable aprox` to install

[Unit]
Description=AProx Local Service
After=network.target

[Service]
Restart=always
RestartSec=10
ExecStart=/opt/aprox/bin/aprox.sh
ExecStop=/usr/bin/pkill -f aprox

[Install]
WantedBy=multi-user.target
